# backend/deb/resources/lib/systemd/system/alfio-configure-printer@.service
[Unit]
Description=Configure Alf.io Printer
Requires=cups.socket
After=cups.socket

[Service]
ExecStart=@ALFIO_UTILS_DIR@/init-printer.sh "%i"